# ops/systemd/personal-data-warehouse-agent-sessions-upload.timer
[Unit]
Description=Run the agent-sessions uploader every 5 minutes

[Timer]
# Mirrors the macOS LaunchAgent cadence (StartInterval=300, RunAtLoad).
OnBootSec=2min
OnUnitActiveSec=300s
# Catch up after the VM was suspended/offline instead of waiting a full period.
Persistent=true
Unit=personal-data-warehouse-agent-sessions-upload.service

[Install]
WantedBy=timers.target
